# packaging/thinkpad-hpd-agent.service
[Unit]
Description=ThinkPad Human Presence Desktop Agent
After=graphical-session.target
PartOf=graphical-session.target
Requires=graphical-session.target

[Service]
ExecStart=/usr/local/bin/thinkpad-hpd agent
Restart=on-failure
RestartSec=3
Environment=XDG_RUNTIME_DIR=%t
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=%t/bus
NoNewPrivileges=true
CapabilityBoundingSet=
PrivateNetwork=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=read-only
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectProc=invisible
ProtectSystem=strict
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
LockPersonality=true
MemoryDenyWriteExecute=true
UMask=0077

[Install]
WantedBy=graphical-session.target
